# components/core/supervisor_api.conf
[program:bassa_api]
directory=/bassa
command=gunicorn -c wsgi_conf.py Main:server
autostart=true
autorestart=true
startretries=5